# supervisor/ngs-oplatform.cfg
[program:ngs-oplatform]
command = /opt/ngs/application/ngs-oplatform/ngs-oplatform --yaml /opt/ngs/application/ngs-oplatform/configs/system.yaml:uwsgi
process_name = %(program_name)s
autostart = true
autorestart = true
user = root
exitcodes = 0
environment = PYTHONPATH="/opt/ngs/application/ngs-oplatform/"
stopasgroup = true
killasgroup = true
stdout_logfile = /opt/ngs/log/ngs-oplatform/stdout.log
stderr_logfile = /opt/ngs/log/ngs-oplatform/stderr.log
stdout_logfile_maxbytes = 20MB  ; stdout 日志文件大小，默认 50MB
stdout_logfile_backups = 20     ; stdout 日志文件备份数
stderr_logfile_maxbytes = 20MB  ; stderr 日志文件大小，默认 50MB
stderr_logfile_backups = 20     ; stderr 日志文件备份数